# systemd/streamdeck-go-helper.service
[Unit]
Description=Stream Deck privileged command helper
Documentation=https://github.com/WoodardDigital/streamdeck-go
# Start before the user session so the socket is ready when streamdeck-go starts.
Before=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/streamdeck-helper
Restart=on-failure
RestartSec=3

# Run as root (required to execute privileged commands).
# The socket is group-restricted to the 'streamdeck' group — only
# members of that group can connect.
RuntimeDirectory=streamdeck-go
RuntimeDirectoryMode=0755

[Install]
WantedBy=multi-user.target
